# service-ttrss-daemon.conf
[program:ttrss-daemon]
command=/usr/bin/php /var/www/ttrss/update_daemon2.php
startsecs=0
startretries=99
user=www-data
redirect_stderr=true
